# system/etc/profile.d/rubygems1.8.sh
#!/bin/bash
# Ruby Gems local
# Make sure it is sudo chmod a+x /etc/profile.d/rubygems1.8.sh
# See: http://groups.google.com/group/compass-users/browse_thread/thread/fc72946d45bf7db5?pli=1
PATH="$PATH:/var/lib/gems/1.8/bin:$HOME/.gem/ruby/1.8/bin"
