# system/etc/profile.d/svn-meld.sh
#!/bin/bash
alias sd='svn diff --diff-cmd /usr/bin/svn-meld'
